# tests/test_seabee.service
[Unit]
Description=SeaBee Service

[Service]
Type=exec
ExecStart=/usr/sbin/seabee --sigint allow --log-level debug
KillSignal=SIGINT
