# install/journald.service
[Unit]
Description=SeaBee Service
RequiresMountsFor=/sys/fs/bpf
After=sysinit.target

[Service]
Type=exec
ExecStart=/usr/sbin/seabee
StandardOutput=null
KillSignal=SIGINT

[Install]
WantedBy=basic.target
